# TP3/systemdunits/backup.timer
[Unit]
Description=Sauvegarde toutes les heures

[Timer]
OnCalendar=*-*-* *:00:00
Unit=lebackup.service

[Install]
WantedBy=multi-user.target
